# scripts/systemd/tailscale-healthcheck.timer
[Unit]
Description=Run TinyWeb Tailscale health verification periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=tailscale-healthcheck@magicdns.service

[Install]
WantedBy=timers.target

