# scripts/systemd/tailscale-healthcheck@.service
[Unit]
Description=TinyWeb Tailscale health verification (%i)
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=oneshot
Environment=TAILSCALE_BIN=/usr/bin/tailscale
Environment=TAILSCALE_HEALTHCHECK_HOST=%i
ExecStart=/usr/bin/env bash /home/william/Documents/codingProj/tinyweb/CTinyWeb/scripts/tailscale_health_check.sh --timeout 5

[Install]
WantedBy=multi-user.target

